# etc/systemd/hive-metastore.service.template
[Unit]
Description=Apache Hive Metastore
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=hive
Group=hive
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64
Environment=HADOOP_HOME=/opt/hadoop
Environment=HIVE_HOME=/opt/hive
ExecStart=/opt/hive/bin/hive --service metastore
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hive-metastore

[Install]
WantedBy=multi-user.target
