# etc/systemd/kafka.service.template
[Unit]
Description=Apache Kafka (single-node KRaft)
After=network.target

[Service]
Type=simple
User=kafka
Group=kafka
ExecStart=/usr/local/bin/kafka-start
Restart=always
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
